# templates/nzbget.upstart.erb
author "Alejandro Figueroa"
description "NZBGet upstart script"

exec <%= "#{@install_dir}/nzbget -D -c #{@config_file}" %>

expect fork
respawn

setuid <%= @user %>

start on runlevel [2345]
stop on runlevel [016]

post-stop script
exec <%= "#{@install_dir}/nzbget -Q" %>
end script
